# build/supervisord.conf
[supervisord]
nodaemon=true

[program:alpino-parser]
command=make -f /Alpino/Makefile.start_server
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true

[program:alpino-tokenizer]
command=/Alpino/Tokenization/tokenize_server/start.sh
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
auto_start=true
autorestart=true